# vpn/nebula/client/nebula.service
[Unit]
Description=Nebula 

[Service]
ExecStart=/etc/nebula/nebula -config /etc/nebula/config.yaml

[Install]
WantedBy=multi-user.target
